# tools/upstart/apertium-all.conf
description "start all apertium services"
     
start on startup
     
script
    initctl emit apertium_start
end script

echo "`date` $UPSTART_JOB started/stopped" >> $LOGFILE
