# tools/upstart/apertium-all.conf
description "start all apertium services"
     
start on startup
